# yokwe-finance/data/crontab
#
# expect macos
#
SHELL=/opt/local/bin/bash
MAILTO=hasegawa
PATH=/opt/local/sbin:/opt/local/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin
# Use english language
LANG=en_US.UTF-8
LANGUAGE=en_US:en

YOKWE_FINANCE=/Users/hasegawa/git/yokwe-root/yokwe-finance


#m   h     dom mon dow   command
# CLEAR SAVE FILE and REPORT FILE
02   1     *   *   *     cd $YOKWE_FINANCE; ( make LOG_TITLE="maintain save file"  log-start clear-save-file clear-report-file log-stop ) 1>>tmp/cron.log 2>&1
# UPDATE FUND
05   3     *   *   6     cd $YOKWE_FINANCE; ( make LOG_TITLE="update trading fund"  log-start update-trading-fund log-stop ) 1>>tmp/cron.log 2>&1
15   3     *   *   1-6   cd $YOKWE_FINANCE; ( make LOG_TITLE="update fund"  log-start update-fund update-nisa report-fund-stats log-stop ) 1>>tmp/cron.log 2>&1
# UPDATE JP STOCK INFO, FUND and DIV
05   4     *   *   1-6   cd $YOKWE_FINANCE; ( make LOG_TITLE="update jp stock info and div"  log-start update-stock-info-jp update-manebu update-jreit update-stock-div-jp log-stop ) 1>>tmp/cron.log 2>&1
# UPDATE JP STOCK PRICE and REPORT
35   4     *   *   1-6   cd $YOKWE_FINANCE; ( make LOG_TITLE="update jp stock price"  log-start update-stock-price-jp report-stock-stats-jp log-stop ) 1>>tmp/cron.log 2>&1
# SAVE FILE
05   5     *   *   *     cd $YOKWE_FINANCE; ( make LOG_TITLE="save file"  log-start save-all log-stop ) 1>>tmp/cron.log 2>&1
# UPDATE ASSET
#15   9     *   *   *     cd $YOKWE_FINANCE; ( make LOG_TITLE="update asset" log-start update-asset-all log-stop  ) 1>>tmp/cron.log 2>&1
#20   9     *   *   *     cd $YOKWE_FINANCE; ( make LOG_TITLE="update asset" log-start update-asset-all report-asset-stats log-stop  ) 1>>tmp/cron.log 2>&1
# UPDATE JP STOCK PRICE and REPORT
55   9     *   *   1-5   cd $YOKWE_FINANCE; ( make LOG_TITLE="update jp stock price"  log-start update-stock-price-jp report-stock-stats-jp log-stop ) 1>>tmp/cron.log 2>&1
# UPDATE JP STOCK PRICE and REPORT
15   12    *   *   1-5   cd $YOKWE_FINANCE; ( make LOG_TITLE="update jp stock price"  log-start update-stock-price-jp report-stock-stats-jp log-stop ) 1>>tmp/cron.log 2>&1
# SAVE FILE
30   12    *   *   *     cd $YOKWE_FINANCE; ( make LOG_TITLE="save file"  log-start save-all log-stop ) 1>>tmp/cron.log 2>&1
# UPDATE US STOCK INFO, DIV, PRICE and REPORT
35   12    *   *   6     cd $YOKWE_FINANCE; ( make LOG_TITLE="update us stock"  log-start update-trading-stock log-stop ) 1>>tmp/cron.log 2>&1
45   12    *   *   2-6   cd $YOKWE_FINANCE; ( make LOG_TITLE="update us stock"  log-start update-stock-info-us update-stock-div-us update-stock-price-us report-stock-stats-us log-stop ) 1>>tmp/cron.log 2>&1
# UPDATE JP STOCK PRICE and REPORT
15   16    *   *   1-5   cd $YOKWE_FINANCE; ( make LOG_TITLE="update jp stock price"  log-start update-stock-price-jp report-stock-stats-jp log-stop ) 1>>tmp/cron.log 2>&1
# UPDATE ASSET
#35   16    *   *   *     cd $YOKWE_FINANCE; ( make LOG_TITLE="update asset" log-start update-asset-all log-stop  ) 1>>tmp/cron.log 2>&1
#40   16    *   *   *     cd $YOKWE_FINANCE; ( make LOG_TITLE="update asset" log-start update-asset-all report-asset-stats log-stop  ) 1>>tmp/cron.log 2>&1
# SAVE FILE
05   23    *   *   *     cd $YOKWE_FINANCE; ( make LOG_TITLE="save file"  log-start save-all log-stop ) 1>>tmp/cron.log 2>&1
# CLEAR LOG FILE
59   23    *   *   6     cd $YOKWE_FINANCE; ( make LOG_TITLE="maintain log file"  log-start clear-log-file log-stop ) 1>>tmp/cron.log 2>&1

# UPDATE FX RATE
10  9,10,15  *  *  1-5   cd $YOKWE_FINANCE; ( make LOG_TITLE="update fx rate"  log-start update-fx-rate log-stop ) 1>>tmp/cron.log 2>&1


# update-stock-us-div takes 48 minutes first time. second time is 9 minutes
# update-stock-us-price takes about 40 minutes
