# yokwe-finance/data/crontab
#
# expect macos
#
SHELL=/opt/local/bin/bash
MAILTO=hasegawa
PATH=/opt/local/sbin:/opt/local/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin
# Use english language
LANG=en_US.UTF-8
LANGUAGE=en_US:en

YOKWE_FINANCE=/Users/hasegawa/git/yokwe-root/yokwe-finance


#m   h     dom mon dow   command
# CLEAR SAVE FILE and REPORT FILE
02   1     *   *   *     cd $YOKWE_FINANCE; ( make clear-save-file clear-report-file ) 1>>tmp/cron.log 2>&1
# UPDATE JP STOCK INFO, FUND and DIV
05   4     *   *   1-5   cd $YOKWE_FINANCE; ( make update-jpx-listing update-stock-info-jp update-manebu update-jreit update-fund update-stock-div-jp ) 1>>tmp/cron.log 2>&1
# SAVE FILE
05   5     *   *   *     cd $YOKWE_FINANCE; ( make save-all ) 1>>tmp/cron.log 2>&1
# UPDATE JP STOCK PRICE and REPORT
55   9     *   *   1-5   cd $YOKWE_FINANCE; ( make update-stock-price-jp report-stock-stats-jp ) 1>>tmp/cron.log 2>&1
# UPDATE JP STOCK PRICE and REPORT
15   12    *   *   1-5   cd $YOKWE_FINANCE; ( make update-stock-price-jp report-stock-stats-jp ) 1>>tmp/cron.log 2>&1
# SAVE FILE
20   12    *   *   *     cd $YOKWE_FINANCE; ( make save-all ) 1>>tmp/cron.log 2>&1
# UPDATE US STOCK INFO, DIV, PRICE and REPORT
30   12    *   *   2-6   cd $YOKWE_FINANCE; ( make update-stock-info-us update-stock-div-us update-stock-price-us report-stock-stats-us ) 1>>tmp/cron.log 2>&1
# UPDATE JP STOCK PRICE and REPORT
45   15    *   *   1-5   cd $YOKWE_FINANCE; ( make update-stock-price-jp report-stock-stats-jp ) 1>>tmp/cron.log 2>&1
# SAVE FILE
05   23    *   *   *     cd $YOKWE_FINANCE; ( make save-all ) 1>>tmp/cron.log 2>&1
# CLEAR LOG FILE
59   23    *   *   6     cd $YOKWE_FINANCE; ( make clear-log-file ) 1>>tmp/cron.log 2>&1


# update-stock-us-div takes 48 minutes first time. second time is 9 minutes
# update-stock-us-price takes about 40 minutes
